# config_lirc/lircd
#!/sbin/openrc-run
# Custom OpenRC service for LIRC

name=$RC_SVCNAME
description="LIRC Daemon for Infrared Remote Control"
command="/usr/sbin/lircd"
pidfile="/var/run/lirc/lircd.pid"
